# HOST/code/capturebox-lora.service
[Unit]
Description=CaptureBox LoRa Receiver
After=network.target
# Ensure Flask app (and DB) is already running
Wants=capturebox.service

[Service]
WorkingDirectory=/home/pi/capturebox
ExecStart=/home/pi/capturebox/venv/bin/python lora_receiver.py
Restart=always
RestartSec=5
User=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
